# deploy/telegram-suite.service
[Unit]
Description=Telegram Admin Suite (main app + collector in one process)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-admin-bot
ExecStart=/opt/telegram-admin-bot/.venv/bin/python -m src.run_all
Restart=always
RestartSec=5
User=telegrambot
Group=telegrambot
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8

[Install]
WantedBy=multi-user.target
