# deploy/telegram-admin-bot.service
[Unit]
Description=Telegram Admin Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-admin-bot
ExecStart=/opt/telegram-admin-bot/.venv/bin/python -m src.main
Restart=always
RestartSec=5
User=telegrambot
Group=telegrambot
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8

[Install]
WantedBy=multi-user.target
